# plists/com.bugclaw.chromium-fuzz-angle-texture-vk-pitch-narrow.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bugclaw.chromium-fuzz-angle-texture-vk-pitch-narrow</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/bugclaw/.openclaw/workspace/chromium-vrp/fuzz/run_angle_texture_vk_pitch_narrow_fuzz.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/bugclaw</string>
        <key>CHROMIUM_FUZZ_LANE_REGISTRY</key>
        <string>/Users/bugclaw/.openclaw/workspace/chromium-vrp/fuzz/managed/lane_registry.json</string>
        <key>MAX_TOTAL_TIME</key>
        <string>172800</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/bugclaw/.openclaw/workspace/chromium-vrp</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/bugclaw/.openclaw/workspace/chromium-vrp/fuzz/logs/com.bugclaw.chromium-fuzz-angle-texture-vk-pitch-narrow.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bugclaw/.openclaw/workspace/chromium-vrp/fuzz/logs/com.bugclaw.chromium-fuzz-angle-texture-vk-pitch-narrow.stderr.log</string>
</dict>
</plist>
