# m1-worker/launchagents/com.bugclaw.chromium-fuzz-media-h264.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bugclaw.chromium-fuzz-media-h264</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/smlee/chromium-worker-ops/bin/lane-loop.sh</string>
        <string>media_h264_decoder_fuzzer</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>60</integer>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>/Users/smlee/chromium-fuzz-data/logs/media-h264-loop.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/smlee/chromium-fuzz-data/logs/media-h264-loop.log</string>
</dict>
</plist>
