# m1-worker/launchagents/com.bugclaw.chromium-worker-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bugclaw.chromium-worker-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/smlee/chromium-worker-ops/bin/sync-repo.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>/Users/smlee/chromium-fuzz-data/logs/worker-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/smlee/chromium-fuzz-data/logs/worker-sync.log</string>
</dict>
</plist>
